# systemd/pev-dashboard.service
[Unit]
Description=PEV Dashboard - Planner Executor Verifier cycle control
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/PEV-scaffold/dashboard
EnvironmentFile=-/home/pi/PEV-scaffold/config/hermes.env
ExecStart=/usr/bin/python3 /home/pi/PEV-scaffold/dashboard/server.py
# Agent turns (headless `claude -p` / `codex exec`) are spawned as transitive
# children of this unit and stay in its cgroup — setsid detaches the session,
# not the cgroup. With the default KillMode=control-group, restarting this
# service SIGKILLs in-flight agent turns. KillMode=process kills only the main
# process, so long-running turns survive a restart.
KillMode=process
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
